# systemd/pi-capture.service
[Unit]
Description=Birdcam Pi Capture
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/birdcam
Environment="PATH=/home/pi/birdcam/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/home/pi/birdcam"
ExecStart=/home/pi/birdcam/.venv/bin/python -m pi_capture.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-capture

# Security settings - relaxed for camera access
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/bird_footage /home/pi/birdcam/database /home/pi/birdcam

# Device access for cameras - allow all video and media devices
PrivateDevices=no
DevicePolicy=auto
SupplementaryGroups=video

# Additional permissions for Picamera2
# Needs access to various subsystem devices
PrivateUsers=no
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target